# vhob-host.service.in
[Unit]
Description=Open URLs in VMs on host 

[Service]
Type=simple
ExecStart=@INSTALL_DIR@/vhob-host-service

[Install]
WantedBy=default.target